# show-nfs.timer
[Unit]
Description=execute showmount in nfs_server

[Timer]
OnCalendar=0:00
Unit=show-nfs.service

[Install]
WantedBy=timers.target

